# packaging/systemd/nbackup-server.service
[Unit]
Description=NBackup Server — Streaming Backup Receiver
Documentation=man:nbackup-server(1) https://github.com/nishisan-dev/n-backup
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/nbackup-server --config /etc/nbackup/server.yaml
Restart=on-failure
RestartSec=5s

# Usuário de serviço
User=nbackup
Group=nbackup

# Limite de file descriptors (uma conexão por agent simultâneo)
LimitNOFILE=65536

# Hardening de segurança
ProtectSystem=strict
NoNewPrivileges=yes
PrivateTmp=yes
ProtectHome=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictSUIDSGID=yes
RestrictNamespaces=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
RestrictRealtime=yes

# O server precisa escrever nos diretórios de storage
ReadWritePaths=/var/backups

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=nbackup-server

[Install]
WantedBy=multi-user.target
