# packaging/systemd/nbackup-agent.service
[Unit]
Description=NBackup Agent — Streaming Backup Client
Documentation=man:nbackup-agent(1) https://github.com/nishisan-dev/n-backup
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/nbackup-agent --config /etc/nbackup/agent.yaml
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=10s

# O agent roda como root para ter acesso de leitura a qualquer path do filesystem.
# Hardening mantido onde possível sem restringir leitura.
NoNewPrivileges=yes
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictSUIDSGID=yes
RestrictNamespaces=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
RestrictRealtime=yes

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=nbackup-agent

[Install]
WantedBy=multi-user.target
