# packaging/systemd/nbackup-server.service
[Unit]
Description=NBackup Server — Streaming Backup Receiver
Documentation=man:nbackup-server(1) https://github.com/nishisan-dev/n-backup
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/nbackup-server --config /etc/nbackup/server.yaml
Restart=on-failure
RestartSec=5s

# Limite de file descriptors (uma conexão por agent simultâneo)
LimitNOFILE=65536

# O server roda como root para escrever em qualquer path de storage configurado.
# Hardening kernel mantido.
NoNewPrivileges=yes
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictSUIDSGID=yes
RestrictNamespaces=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
RestrictRealtime=yes

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=nbackup-server

[Install]
WantedBy=multi-user.target
