# packaging/systemd/nbackup-agent.service
[Unit]
Description=NBackup Agent — Streaming Backup Client
Documentation=man:nbackup-agent(1) https://github.com/nishisan-dev/n-backup
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/nbackup-agent --config /etc/nbackup/agent.yaml
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=10s
LimitNOFILE=65536
Nice=0
IOSchedulingClass=best-effort
IOSchedulingPriority=4

StandardOutput=journal
StandardError=journal
SyslogIdentifier=nbackup-agent

[Install]
WantedBy=multi-user.target
