# packaging/systemd/nbackup-server.service
[Unit]
Description=NBackup Server — Streaming Backup Receiver
Documentation=man:nbackup-server(1) https://github.com/nishisan-dev/n-backup
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/nbackup-server --config /etc/nbackup/server.yaml
Restart=on-failure
RestartSec=5s
LimitNOFILE=65536
Nice=0
IOSchedulingClass=best-effort
IOSchedulingPriority=4

StandardOutput=journal
StandardError=journal
SyslogIdentifier=nbackup-server

[Install]
WantedBy=multi-user.target
